# debug-env/roles/wpa_supplicant.debug/templates/wpa_supplicant.service.d/10-exec-start.conf.j2
[Service]
ExecStart=
ExecStart=/usr/local/bin/wpa_supplicant.debug -u -s -O /run/wpa_supplicant
